# aahwebsite.service
# GENERATED BY aah CLI - Feel free to customization it.
# FILE: aahwebsite.service
# DATE: Fri, 07 Dec 2018 23:40:48 -0800
# DESC: aah application systemd service file

[Unit]
Description=aahwebsite application
After=network.target

[Service]
User=aah
Group=aah
EnvironmentFile=/home/aah/aahwebsite_env_values
ExecStart=/home/aah/aahwebsite/bin/aahwebsite run --envprofile prod
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
