# aahwebsite.service
# GENERATED BY aah CLI - Feel free to customization it.
# FILE: aahwebsite.service
# DATE: Sat, 07 Jul 2018 00:25:34 -0700
# DESC: aah application systemd service file

[Unit]
Description=aahwebsite application
After=network.target

[Service]
User=aah
Group=aah
EnvironmentFile=/home/aah/aahwebsite_env_values
ExecStart=/home/aah/aahwebsite/bin/aahwebsite -profile prod
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
